# deploy/prod/rbp-functor.rpc.voice_lover_admin.conf
[program:rbp-functor.rpc.voice_lover_admin]
directory=/home/ecs-user/webroot/rbp-functor.rpc.voice_lover_admin
command=/home/ecs-user/webroot/rbp-functor.rpc.voice_lover_admin/bin/rpc --name=voice_lover_admin --gf.gcfg.file=config_prod.toml
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/home/ecs-user/log/rbp-functor.rpc.voice_lover_admin/rbp-functor.rpc.voice_lover_admin.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stdout_capture_maxbytes=1MB
redirect_stderr=true
stderr_logfile=/home/ecs-user/log/rbp-functor.rpc.voice_lover_admin/rbp-functor.rpc.voice_lover_admin.stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
environment = APP_VERSION="8"
